# templates/service/podium-registry.service.erb
[Unit]
Description=Podium Registry microservice

[Service]
User=<%= @user %>
WorkingDirectory=<%= @home %>
ExecStart=<%= @registry_start_script %>
Restart=always
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target

