# templates/service/podium-uaa.service.erb
[Unit]
Description=Podium User Authentication microservice
After=podium-registry.service network.target

[Service]
User=<%= @user %>
WorkingDirectory=<%= @home %>
ExecStartPre=/bin/sleep 10
ExecStart=<%= @uaa_start_script %>
Restart=always
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target

